# client/gcast-client.service.in
[Unit]
Description=web server for gcast

[Service]
ExecStart=/bin/python ./serve.py
WorkingDirectory=%h/${IN_DEPLOY_PATH}/client
