# server/gcast-server.service.in
[Unit]
Description=server for gcast

[Service]
ExecStart=%h/${IN_DEPLOY_PATH}/server/server
WorkingDirectory=%h/${IN_DEPLOY_PATH}/server
Environment=RUST_BACKTRACE=full
